# build/package/aerospike-backup-service.service
[Unit]
Description=Aerospike Backup Service
Wants=network.target
After=network-online.target

[Service]
ExecStart=/usr/local/bin/aerospike-backup-service -c /etc/aerospike-backup-service/config.yml

[Install]
WantedBy=default.target
