# build/package/post-install.sh
#!/bin/bash
systemctl enable aerospike-backup-service.service
systemctl daemon-reload
