# packages/debian/aerospike-backup-service.service
[Unit]
Description=Aerospike Backup Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/aerospike-backup-service --config /etc/aerospike-backup-service/config.yml
Restart=on-failure

[Install]
WantedBy=multi-user.target