# build/package/scripts/preremove.sh
#!/bin/bash -e
systemctl stop aerospike-backup-service
systemctl disable aerospike-backup-service
